# deploy_tools/rabbitmq.template.conf
[program:rabbitmq]
command=/usr/sbin/rabbitmq-server
user=root
autostart=true
autorestart=true
priority=10
redirect_stderr=true
stdout_logfile=/home/tiago/sites/SITENAME/rabbitmq-server.log
